# archive/bsse/modules/profile/templates/rc.local.erb
#
# This file was autogenerated on <%= File.mtime("#{file}").strftime("%m/%d/%Y %T %Z") %> # by puppet <%= puppetversion %>
# It can still be managed manually but will be overwritten on the next run.
# SVN <%= %x(svn info #{file} | grep Revision | tr -d "\n")%>
#

<% if classes.include?("service-iscsi-client") %>
#
#	Assuming eth4/eth5 are the iSCSI networks
#
ethtool --offload eth4 gro off
ethtool --offload eth5 gro off

#	Determine if the standard networks 129.132.27.0 or 195.176.122.0 are used for iSCSI and if so, disable them
if [ -e "/usr/sbin/rswcli" ]; then
	if [ ! -z "`rswcli -L | sed ':a;N;$!ba;s/Subnets.excluded.*//' | grep -v "^$" | grep "129.132.27"`" ]; then
		rswcli -E --net 129.132.27.0 --mask 255.255.255.128
	fi
	if [ ! -z "`rswcli -L | sed ':a;N;$!ba;s/Subnets.excluded.*//' | grep -v "^$" | grep "195.176.122"`" ]; then
		rswcli -E --net 195.176.122.0 --mask 255.255.255.0
	fi
fi

<% end %>
<% if classes.include?("hw-kvm") %>
for i in `echo "/dev/vda /dev/vdb /dev/vdc /dev/vdd"`; do
	if [ -b "$i" ]; then
		blockdev --setra 65536 $i
	fi
done
<% end %>
<% if classes.include?("ou-service-multihomed") %>
#
#	Test if the internal networks are existant, if not create.
#
if [ -z "`grep srvNet /etc/iproute2/rt_tables`" ]; then
	logger "Adding srvNet to iproute2"
	echo 200 srvNet >> /etc/iproute2/rt_tables
fi
if [ -z "`grep dockNet /etc/iproute2/rt_tables`" ]; then
	logger "Adding dockNet to iproute2"
	echo 201 dockNet >> /etc/iproute2/rt_tables
fi
if [ -z "`grep cinaServerNet /etc/iproute2/rt_tables`" ]; then
	logger "Adding cinaServerNet to iproute2"
	echo 202 cinaServerNet >> /etc/iproute2/rt_tables
fi

#
#	add default routes to the according networks
#
<%# 	Very static approach atm... %>
<% if network_eth0 == "195.176.122.0" %>
ip route add default via 195.176.122.1 table srvNet <% else %>
ip route add default via 129.132.27.1 table srvNet 
<% end %>
ip route add default via 129.132.228.1 table dockNet
ip route add default via 129.132.151.1 table cinaServerNet

#
#	add routes to the according networks
#
ip rule add from 129.132.27.0/25 table srvNet
ip rule add from 195.176.122.0/24 table srvNet
ip rule add from 129.132.228.0/23 table dockNet
ip rule add from 129.132.151.0/28 table cinaServerNet
ip route flush cache
<% end %>

#	Run puppet once on startup
/usr/bin/puppet agent -tv

