# archive/bsse/modules/udev/templates/99_vmware_scsi.rules.erb
#
# This file was autogenerated on <%= File.mtime("#{file}").strftime("%m/%d/%Y %T %Z") %> # by puppet <%= @puppetversion %>
# It can still be managed manually but will be overwritten on the next run.
# SVN <%= %x(svn info #{file} | grep Revision | tr -d "\n")%>
#
# SCSI devices Timeout adjustment
#
# Modify the timeout value for SCSI devices so that
# in the event of a failover, we don't time out.
# See Bug 271286 for more information.
#
<% if @virtual =~ /vmware/ && @lsbmajdistrelease == '6' -%>
ACTION=="add", SUBSYSTEM=="scsi" , SYSFS{type}=="0|7|14", RUN+="/bin/sh -c 'echo 180 > /sys$$DEVPATH/timeout'"
<% elsif @virtual =~ /vmware/ && @lsbmajdistrelease == '7' -%>
ACTION=="add", SUBSYSTEMS=="scsi", RUN+="/bin/sh -c 'echo 180 >/sys$DEVPATH/timeout'"
<% end -%>
